# timers/docker.timer
[Unit]
Description="Delayed start for the Docker service"

[Timer]
OnBootSec=30s

[Install]
WantedBy=timers.target
